# speedtest.cron
7 * * * * /bin/python3 /home/scripts/speedtest.py  --simple >> /var/log/speed